# wpa_supplicant/wpa_supplicant.service
[Unit]
Description=Wi-Fi Security Service

[Service]
Type=dbus
BusName=fi.epitest.hostap.WPASupplicant
ExecStart=/usr/sbin/wpa_supplicant -c /etc/wpa_supplicant/wpa_supplicant.conf -u

# Hmm, a WantedBy=multi-user.target probably makes sense here.
